# iodbc-dsn-manager/server/com.openlink.odbc-rest-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openlink.odbc-rest-server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/usr/local/share/odbc-rest-server/odbc_rest_server.py</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8899</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/odbc-rest-server.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/odbc-rest-server.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/share/odbc-rest-server</string>
</dict>
</plist>
